# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" -c /app/nginx.conf
autostart=true
autorestart=true
startretries=5
priority=10
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err.log

[program:backend]
command=node /app/backend/server.js
directory=/app/backend
autostart=true
autorestart=true
startretries=10
priority=20
environment=NODE_ENV="production",API_PORT="3001"
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend.err.log
